# deploy/systemd/bistbot.service
[Unit]
Description=Telegram BIST Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=bistbot
Group=bistbot
WorkingDirectory=/opt/telegram-bist-bot
Environment=ENV_FILE=/opt/telegram-bist-bot/.env
ExecStart=/opt/telegram-bist-bot/bistbot
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
